# detect/testdata/sigma_rules/real_world/persistence_scheduled_task.yml
title: Suspicious Scheduled Task Creation
id: sigma-real-005
status: test
description: Detects scheduled task creation for persistence
logsource:
    category: process_creation
    product: windows
detection:
    selection_schtasks:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains: '/create'
    selection_persistence:
        CommandLine|contains:
            - '/sc onstart'
            - '/sc onlogon'
            - '/sc daily'
            - '/rl highest'
    filter_system:
        User: 'SYSTEM'
        CommandLine|contains:
            - 'Microsoft\Windows'
            - 'WindowsUpdate'
    condition: selection_schtasks and selection_persistence and not filter_system
level: medium
tags:
    - attack.persistence
    - attack.t1053.005
